# sources/meta-user/recipes-apps/portfolio-website/portfolio-website/portfolio-website.service
[Unit]
Description=Startup Portfolio Server

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/portfolio-website/start-portfolio-website
